# agent/systemd/cloudwatch-agent.timer
[Unit]
Description=Regularly send RAM usage metrics via metadata API

[Timer]
OnCalendar=*-*-* *:*:0
AccuracySec=1

[Install]
WantedBy=timers.target
