# DEBIAN/prerm
#!/bin/bash

systemctl stop cloudwatch-agent.timer
systemctl disable cloudwatch-agent.timer

rm /usr/lib/systemd/system/cloudwatch-agent.service
rm /usr/lib/systemd/system/cloudwatch-agent.timer
